# infrastructure/terraform/madoc/single-instance/files/backup/madoc-db-backup-daily.timer
[Unit]
Description=MySQL daily backup

[Timer]
OnCalendar=02:30
Persistent=true
Unit=madoc-db-backup@daily.service

[Install]
WantedBy=timers.target