# infrastructure/terraform/madoc/single-instance/files/backup/madoc-backup.service
[Unit]
Description=Backup madoc files and mysql db

[Service]
Type=simple
EnvironmentFile=/etc/madoc-environment
ExecStart=/usr/bin/backup-madoc.sh

[Install]
WantedBy=multi-user.target